# roles/apps/templates/electron_rasp-update.service.j2
[Unit]
Description=Pull electron_rasp repo y reinicia servicios si cambia
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ kiosk_user }}
WorkingDirectory={{ electron_root }}/electron_rasp
Environment=HOME=/home/{{ kiosk_user }}
ExecStart={{ electron_root }}/update_repo.sh
Restart=no
StandardOutput=journal
StandardError=journal