# roles/apps/templates/electron_rasp-electron.service.j2
[Unit]
Description=Electron kiosk (electron_rasp)
After=graphical.target electron_rasp-flask.service
Requires=electron_rasp-flask.service
StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
User={{ kiosk_user }}
WorkingDirectory={{ electron_app_dir }}

Environment=DISPLAY=:0
Environment="BRANCH_ID={{ BRANCH_ID }}"
Environment="KIOSK_URL=http://127.0.0.1:5000"
Environment="SPAWN_FLASK=0"
ExecStart=/usr/bin/npm run electron

KillMode=control-group
TimeoutStopSec=15s
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target