# roles/apps/templates/signage-update.service.j2
[Unit]
Description=Pull signage repo y reinicia electron_rasp si cambia
After=network-online.target
Wants=network-online.target
ConditionPathExists={{ signage_root }}/signage/.git

[Service]
Type=oneshot
User={{ kiosk_user }}
Group={{ kiosk_user }}
WorkingDirectory={{ signage_root }}/signage

# Entorno correcto + fuerza llave de signage (cinturón y tirantes)
Environment=HOME=/home/{{ kiosk_user }}
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=GIT_SSH_COMMAND=/usr/bin/ssh -o IdentitiesOnly=yes -i /home/{{ kiosk_user }}/.ssh/id_ed25519_segundo -o StrictHostKeyChecking=accept-new

ExecStart={{ signage_root }}/signage/update.sh
StandardOutput=journal
StandardError=journal

PrivateTmp=yes
ProtectSystem=full
ProtectHome=no
NoNewPrivileges=yes
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target