# roles/apps/templates/signage-update.service.j2
[Unit]
Description=Pull signage repo y reinicia electron_rasp si cambia
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/signage/signage/.git

[Service]
Type=oneshot
User=rasp-trolley
Group=rasp-trolley
WorkingDirectory=/opt/signage/signage

# 🔒 Asegura entorno correcto para git/ssh
Environment=HOME=/home/rasp-trolley
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=GIT_SSH_COMMAND=/usr/bin/ssh -o IdentitiesOnly=yes -i /home/rasp-trolley/.ssh/id_ed25519_signage -o StrictHostKeyChecking=accept-new

# Si usas variables extra:
# EnvironmentFile=-/etc/default/signage-update

ExecStart=/opt/signage/signage/update.sh
StandardOutput=journal
StandardError=journal

# Seguridad razonable para oneshot
PrivateTmp=yes
ProtectSystem=full
ProtectHome=no
NoNewPrivileges=yes
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target