# roles/maintenance/templates/rpi-maintenance.timer.j2
[Unit]
Description=Run Raspberry Pi maintenance daily

[Timer]
OnBootSec=3min
OnCalendar=*-*-* {{ maintenance_hour }}
Persistent=yes

[Install]
WantedBy=timers.target
