# roles/maintenance/templates/rpi-maintenance.service.j2
[Unit]
Description=Raspberry Pi daily maintenance

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rpi_maintenance.sh
SuccessExitStatus=0
