# roles/maintenance/templates/wake-hdmi.timer.j2
[Unit]
Description=Wake HDMI/TV after nightly maintenance

[Timer]
OnCalendar=*-*-* 03:02:00
Persistent=true

[Install]
WantedBy=timers.target